# scripts/ndi-controller.service
[Unit]
Description=NDI Controller Web Interface
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ndi
Group=ndi
WorkingDirectory=/opt/ndi-controller/backend

# Core environment
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=PYTHONUNBUFFERED=1

# Display access for SDL2 window output (yuri2 viewer)
# These are set by the installer based on the desktop user
Environment=DISPLAY=:0
Environment=WAYLAND_DISPLAY=wayland-0
Environment=SDL_VIDEODRIVER=x11
Environment=XDG_RUNTIME_DIR=/run/user/1000
# XAUTHORITY is set dynamically by the installer

# Library paths for yuri and NDI
Environment=LD_LIBRARY_PATH=/usr/local/lib
Environment=NDI_PATH=/usr/local/lib/libndi.so.6
Environment=NDI_RUNTIME_DIR_V6=/usr/share/ndi

# Application paths
Environment=YURI_BIN=/usr/local/bin/yuri2
Environment=YURI_LIB_PATH=/usr/local/lib
Environment=NDI_LIB_PATH=/usr/local/lib/libndi.so.6
Environment=FRONTEND_DIR=/opt/ndi-controller/frontend/dist
Environment=CONFIG_DIR=/opt/ndi-controller/configs/generated
Environment=NDI_EXTRA_IPS_FILE=/opt/ndi-controller/configs/extra_ips.txt

# Use gunicorn with gevent workers for proper streaming support
# -k gevent: use gevent async workers (handles streaming without blocking)
# -w 1: single worker to maintain yuri process state
# --timeout 300: 5 minute timeout for long operations
ExecStart=/opt/ndi-controller/venv/bin/gunicorn -k gevent -w 1 -b 0.0.0.0:5000 --timeout 300 --access-logfile - --error-logfile - app:app

# Graceful shutdown
ExecStop=/bin/kill -TERM $MAINPID
TimeoutStopSec=10

Restart=always
RestartSec=5

# Allow access to video devices, display, and ramdisk for preview
SupplementaryGroups=video audio render
ReadWritePaths=/dev/shm/extrashot_preview

[Install]
WantedBy=multi-user.target
